# systemd/agent-morning-brief.timer
[Unit]
Description=Developer Workflow Morning Brief Timer
Requires=jira-mcp.service

[Timer]
OnCalendar=Mon-Fri 09:03:00
Persistent=true

[Install]
WantedBy=timers.target
