# systemd/agent-hourly-brief.timer
[Unit]
Description=Developer Workflow Hourly Brief Timer
Requires=jira-mcp.service

[Timer]
OnCalendar=Mon-Fri 11:50:00
OnCalendar=Mon-Fri 15:30:00
Persistent=true

[Install]
WantedBy=timers.target
