# systemd/agent-hourly-brief.service
[Unit]
Description=Developer Workflow Hourly Brief

[Service]
Type=oneshot
WorkingDirectory=%h/git/agent
EnvironmentFile=%h/git/agent/.env
ExecStart=UV_PATH run python main.py --mode brief
StandardOutput=journal
StandardError=journal
